# scp-backup/rootfs/etc/services.d/scp-backup/finish
#!/command/execlineb -P
# Home Assistant addon finish script

# Cleanup on service stop
echo "SCP Backup addon stopping..."
