# scp-backup/rootfs/etc/services.d/scp-backup/run
#!/command/execlineb -P
# Home Assistant addon service script

with-contenv

/usr/local/bin/main.sh
